# mac-file-workflow-builder/com.finder.organizer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.finder.organizer</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>__SCRIPT_PATH__</string>
    <string>--sources</string>
    <string>__HOME__/Downloads</string>
    <string>__HOME__/Desktop</string>
  </array>
  <key>StartInterval</key>
  <integer>1800</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/finder-organizer.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/finder-organizer.err.log</string>
</dict>
</plist>
